# crontab
*/19 * * * * /root/killv3.sh
58 23 * * * /root/rebootv3.sh
